# templates/openvpn.service.erb
[Unit]
Description=OpenVPN <%= @title %>
After=network.target auditd.service

[Service]
ExecStart=/usr/sbin/openvpn --config /etc/openvpn/<%= @title %>.conf

[Install]
WantedBy=multi-user.target
